# system_files/usr/lib/systemd/system/bazzite63-flatpak-manager.service
[Unit]
Description=Install bazzite-63 default Flatpaks
# After Bazzite's manager so flathub is enabled and the blocklist filter is
# already applied when we install.
Wants=network-online.target
After=network-online.target bazzite-flatpak-manager.service
# Unlimited retries by design: a no-network first boot must keep trying.
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/libexec/bazzite63-flatpak-manager
Restart=on-failure
RestartSec=90

[Install]
WantedBy=multi-user.target
